# portal/admin_agent/systemd/prevencion-admin-agent.service.example
# /etc/systemd/system/prevencion-admin-agent.service
# Ajusta: User, Group, WorkingDirectory, EnvironmentFile
[Unit]
Description=Prevención admin_agent (OpenRouter proxy, localhost only)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/prevencion/portal/admin_agent
EnvironmentFile=/var/www/prevencion/portal/admin_agent/.env
ExecStart=/var/www/prevencion/portal/admin_agent/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 9102
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
